# supervisor/mysql.ini
[group:mysql]
programs=server

[program:server]
command=/bin/su mysql /bin/bash -c "/usr/sbin/mysqld --defaults-file=/opt/mysql/conf/node.conf"
process_name=%(program_name)s
pidfile=/opt/mysql/system/mysqld.pid
numprocs=1
autorestart=false
autostart=true
startsecs=10
startretries=1
stopsignal=TERM
stopwaitsecs=10
priority=200
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true